# low-ram-check.service
[Unit]
Description=Low RAM warning

[Service]
ExecStart=/script/path.sh

[Install]
WantedBy=default.target
